# nuki-monitor.service
[Unit]
Description=Nuki Smart Lock Notification System
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/root/nukiweb
ExecStart=/root/nukiweb/venv/bin/python /root/nukiweb/scripts/nuki_monitor.py
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
